# profile.d/_hackaday.sh
function _hackaday {
   COMPREPLY=( "hackaday.com" "hackaday.io" )
   return 0
}
