# profile.d/android_tools.sh
: ${ANDROID_HOME:=$HOME/Library/Android/sdk}
: ${JAVA_HOME:=/Library/Java/JavaVirtualMachines/jdk1.8.0_144.jdk/Contents/Home}
PATH="$ANDROID_HOME/tools:$ANDROID_HOME/tools/bin:$ANDROID_HOME/platform-tools:$JAVA_HOME/bin:$PATH"
export PATH
