# profile.d/ansible1.9

[ -d "/usr/local/opt/ansible@1.9/bin" ] && export PATH="/usr/local/opt/ansible@1.9/bin:$PATH"
