# profile.d/genyshell-aliases.sh
[ -x "/Applications/Genymotion Shell.app/Contents/MacOS/genyshell" ] && \
  alias geny="/Applications/Genymotion\ Shell.app/Contents/MacOS/genyshell"
