# profile.d/keybase.sh
# keybase aliases

alias kbc='keybase chat'
alias kbcs='keybase chat send'
alias kbcu='keybase chat upload'
alias kbcsks='keybase chat send kickserv_dev'
alias kbcuks='keybase chat upload kickserv_dev'
