# profile.d/python-aliases.sh
alias pe=pipenv
alias pep="pipenv run python"
alias per="pipenv run"
alias pyt="pipenv run python -m pytest"
