# profile.d/sdkman.sh
#!/bin/bash
#THIS MUST BE AT THE END OF THE FILE FOR SDKMAN TO WORK!!!
export SDKMAN_DIR="/Users/tamara/.sdkman"
[[ -s "/Users/tamara/.sdkman/bin/sdkman-init.sh" ]] && source "/Users/tamara/.sdkman/bin/sdkman-init.sh"
