# profile.d/usr-local
#!/bin/sh
# append /usr/local stuff to PATH
export PATH=$PATH:/usr/local/bin:/usr/local/sbin
